# systemd/.config/systemd/user/feh.service
[Unit]
Description=feh

[Service]
Type=oneshot
ExecStart=/usr/bin/sh %h/.fehbg

[Install]
WantedBy=default.target
