# systemd/.config/systemd/user/parcellite.service
[Unit]
Description=Parcellite

[Service]
ExecStart=/usr/bin/parcellite

[Install]
WantedBy=default.target
