# systemd/.config/systemd/user/redshift-gtk.service
[Unit]
Description=Redshift GTK

[Service]
ExecStart=/usr/bin/redshift-gtk

[Install]
WantedBy=default.target
