# systemd/.config/systemd/user/setxkbmap.service
[Unit]
Description=setxkbmap layout and remappings

[Service]
ExecStart=/usr/bin/setxkbmap fr -option caps:escape

[Install]
WantedBy=default.target
